# deploy/latin-vocab.service
[Unit]
Description=Latin Vocab Trainer Flask Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/latin-vocab-trainer
Environment="PATH=/var/www/latin-vocab-trainer/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/var/www/latin-vocab-trainer/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 4 flask_app.app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target